# service/webvps.service
[Unit]
Description=Web VPS Manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/web-vps-manager
ExecStart=/root/web-vps-manager/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 5000 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/root/web-vps-manager

[Install]
WantedBy=multi-user.target
